# kerpad.service.template
[Unit]
Description=Edge motion
After=default.target

[Service]
ExecStart=/usr/bin/kerpad <args>
Type=exec
Restart=always

[Install]
WantedBy=default.target